# systemd/gitremind.timer
[Unit]
Description=Py-GitRemind Timer

[Timer]
OnBootSec=30min
OnCalendar=*:0/30
Unit=gitremind.service

[Install]
WantedBy=timers.target
